# deploy/systemd/agents-on-lightning-monitor.service
[Unit]
Description=Agents on Lightning Monitor Dashboard
After=network-online.target agents-on-lightning.service
Wants=network-online.target

[Service]
Type=simple
User=agentsonlightning
Group=agentsonlightning
WorkingDirectory=/opt/agents_on_lightning
EnvironmentFile=/etc/agents-on-lightning/agents-on-lightning-monitor.env
ExecStart=/usr/bin/node /opt/agents_on_lightning/monitoring_dashboards/server.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agents_on_lightning/data
ReadOnlyPaths=/opt/agents_on_lightning
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
LockPersonality=true
UMask=0077

[Install]
WantedBy=multi-user.target
